# IaC/macos/com.lab.minikube-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd daemon — wraps minikube tunnel with a wait-for-Docker/cluster loop.
  Runs as root (LaunchDaemon default); wrapper script uses 'su - markpadam'
  for Docker/minikube checks so the Colima socket is accessible.

  Install (first time):
    sudo cp IaC/macos/minikube-tunnel.sh /usr/local/bin/minikube-tunnel.sh
    sudo chmod +x /usr/local/bin/minikube-tunnel.sh
    sudo cp IaC/macos/com.lab.minikube-tunnel.plist /Library/LaunchDaemons/
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.lab.minikube-tunnel.plist

  Update (after plist or script change):
    sudo launchctl bootout system/com.lab.minikube-tunnel
    sudo cp IaC/macos/minikube-tunnel.sh /usr/local/bin/minikube-tunnel.sh
    sudo chmod +x /usr/local/bin/minikube-tunnel.sh
    sudo cp IaC/macos/com.lab.minikube-tunnel.plist /Library/LaunchDaemons/
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.lab.minikube-tunnel.plist

  Logs:  /var/log/minikube-tunnel.log

  NOTE: Do NOT add a UserName key — it causes "Bootstrap failed: 5" on
  macOS Ventura+ for LaunchDaemons. The daemon runs as root by default.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lab.minikube-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/minikube-tunnel.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/var/log/minikube-tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/minikube-tunnel.log</string>
</dict>
</plist>
